# docker/s6-rc.d/jiny-m/run
#!/command/execlineb -P
cd /opt/jiny-m
/usr/local/bin/jiny-m monitor --workdir /opt/jiny-m
